# images/router/supervisor_router.conf
[program:frr]
command=/usr/lib/frr/frrinit.sh start
autostart=false
autorestart=true
startsecs=0
minfds=1024
